# deploy/uzum-order.service
[Unit]
Description=Uzum -> MoySklad buyurtma sinxronizatsiyasi (bir martalik ishga tushirish)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# --- Quyidagilarni o'z serveringizga moslang ---
User=uzum
WorkingDirectory=/opt/uzumOrderToMC
ExecStart=/usr/bin/node /opt/uzumOrderToMC/src/index.js
# node yo'lini `which node` bilan aniqlang; nvm ishlatilsa to'liq yo'lni yozing.
