# deploy/uzum-order.timer
[Unit]
Description=Uzum -> MoySklad sinxronizatsiyasini har 2 daqiqada ishga tushiradi

[Timer]
# Serverdan 1 daqiqa o'tib birinchi marta ishga tushadi, keyin har 2 daqiqada.
OnBootSec=1min
OnUnitActiveSec=2min
# Agar oldingi ishga tushirish hali tugamagan bo'lsa, yangi nusxa boshlanmaydi
# (systemd service allaqachon "active" bo'lsa timer uni qayta ishga tushirmaydi).
AccuracySec=30s
Unit=uzum-order.service

[Install]
WantedBy=timers.target
